# bluetooth_2_usb.service
[Unit]
Description=Bluetooth to USB HID proxy
After=multi-user.target

[Service]
User=root
ExecStart=/usr/bin/bluetooth_2_usb --input_devices /dev/input/event1,/dev/input/event2
Environment=PYTHONUNBUFFERED=1
Restart=on-failure

[Install]
WantedBy=multi-user.target